# resources/etc/xdg/autostart/lcarsde-app-menu.desktop
[Desktop Entry]
Type=Application
Name=LCARSDE application menu
Exec=/usr/bin/lcarsde-app-menu.py
OnlyShowIn=lcarsde;
